# debian/hashbash-backend.hashbash-webapp.service
[Unit]
Description=hashbash webapp service

[Service]
ExecStart=/opt/hashbash/bin/hashbash-webapp
Restart=always
KillSignal=SIGQUIT
Type=simple
StandardError=syslog
StandardOutput=/opt/hashbash/logs/webapp.log
NotifyAccess=all
EnvironmentFile=/opt/hashbash/hashbash.env
Environment=HASHBASH_MANAGEMENT_PORT=8081

[Install]
WantedBy=multi-user.target
